# gitsign-credential-cache-bin/gitsign-credential-cache.service
[Unit]
Description=GitSign credential cache
PartOf=gitsign-credential-cache.service

[Service]
Type=simple
ExecStart=/usr/bin/gitsign-credential-cache

[Install]
WantedBy=default.target
